# docker/emiapp.ini
[program:h2]
command = /usr/bin/java -cp /opt/h2/h2-1.4.193.jar org.h2.tools.Server -tcp -tcpPort 10888
startsecs = 5
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/h2/h2.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups = 5

[program:emiapp]
command = /usr/bin/java -Dspring.datasource.url=jdbc:h2:tcp://localhost:10888//var/emiapp/emi -jar /opt/emiapp/emiapp-0.0.1.jar
startsecs = 20
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/emiapp/emiapp.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 5