# docker/emiapp.ini
[program:h2]
command = /usr/bin/java -cp /opt/h2/h2-1.4.193.jar org.h2.tools.Server -tcp -tcpPort 45777
startsecs = 5
autorestart = true
redirect_stderr = false

[program:emiapp]
command = /usr/bin/java -Dspring.datasource.url=jdbc:h2:tcp://localhost:45777//var/emiapp/emi -jar /opt/emiapp/emiapp-0.0.1.jar
startsecs = 20
autorestart = true
redirect_stderr = false

[program:postfix]
command = /usr/sbin/postfix start
startsecs = 5
autorestart = true
redirect_stderr = false
